# rpm/files/systemd/platform-api-server.service
[Unit]
Description=Platform - API Server
After=platform-infra.target
Requires=platform-infra.target
PartOf=platform-all.target

[Service]
Type=exec
User=nobody
Group=nobody
WorkingDirectory=/opt/platform/apps/api-server

# Load shared configurations first
EnvironmentFile=/opt/platform/apps/conf/env.properties
EnvironmentFile=/opt/platform/apps/conf/redis.properties

# Load app-specific configuration (can override shared)
EnvironmentFile=/opt/platform/apps/api-server/api-server.properties

# Optional override file (won't fail if missing)
EnvironmentFile=-/opt/platform/apps/api-server/api-server.override

ExecStart=/opt/platform/apps/api-server/api-server
Restart=always
RestartSec=10
StandardOutput=append:/var/log/platform/api-server/stdout.log
StandardError=append:/var/log/platform/api-server/stderr.log
# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/platform/api-server

[Install]
WantedBy=multi-user.target

