# rpm/files/systemd/platform-checkout-api.service
[Unit]
Description=Platform - Checkout API
After=platform-infra.target
Requires=platform-infra.target
PartOf=platform-all.target

[Service]
Type=exec
User=nobody
Group=nobody
WorkingDirectory=/opt/platform/apps/checkout-api

# Load shared configurations first
EnvironmentFile=/opt/platform/apps/conf/env.properties
EnvironmentFile=/opt/platform/apps/conf/redis.properties

# Load app-specific configuration (can override shared)
EnvironmentFile=/opt/platform/apps/checkout-api/checkout-api.properties

# Optional override file (won't fail if missing)
EnvironmentFile=-/opt/platform/apps/checkout-api/checkout-api.override

ExecStart=/opt/platform/apps/checkout-api/checkout-api
Restart=always
RestartSec=10
StandardOutput=append:/var/log/platform/checkout-api/stdout.log
StandardError=append:/var/log/platform/checkout-api/stderr.log
# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/platform/checkout-api

[Install]
WantedBy=multi-user.target

